# deploy/systemd/openai-proxy.service
[Unit]
Description=OpenAI Reverse Proxy (Rust, rustls) - openai-proxy
Documentation=https://github.com/digkill/openai-proxy
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/openai-proxy
EnvironmentFile=/opt/openai-proxy/.env
ExecStart=/opt/openai-proxy/openai-proxy
User=ctobot
Group=ctobot
Restart=on-failure
RestartSec=2s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=~@mount @swap @module @reboot @raw-io @debug
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

# логирование в journald
StandardOutput=journal
StandardError=journal